# launchd/com.ianburke.narrative-tick.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ianburke.narrative-tick</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/ianburke/Documents/GitHub/narrative-parallax/scripts/tick-runner.sh</string>
    </array>

    <key>StartInterval</key>
    <integer>60</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/ianburke/Documents/GitHub/narrative-parallax/logs/tick/tick.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/ianburke/Documents/GitHub/narrative-parallax/logs/tick/tick_error.log</string>
</dict>
</plist>
